# systemd/user/xdg-desktop-portal.service.d/override.conf
# Fix OBS screen capture on Hyprland:
# xdg-desktop-portal.service requires graphical-session.target which
# Hyprland doesn't activate by default (unlike DM-launched sessions).
[Unit]
Requisite=
Wants=graphical-session.target
